# nfs-volume-plugin/nfs-volume-plugin.service
[Unit]
Wants=multi-user.target
After=multi-user.target

[Service]
ExecStart=/nfs-volume-plugin
EnvironmentFile=/pluginenv
