# centos-mounted-volume-plugin/centos-mounted-volume-plugin.service
[Unit]
Wants=multi-user.target
After=multi-user.target

[Service]
ExecStart=/centos-mounted-volume-plugin
EnvironmentFile=/pluginenv
